# src/main/deb/control/prerm
#!/bin/sh

service @project.artifactId@ stop
systemctl disable @project.artifactId@.service
systemctl daemon-reload
